# launchd/com.goldman.livefeed.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.goldman.livefeed.hourly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/sh</string>
      <string>-c</string>
      <string>LOG_FILE=/Users/eric/.goldman-news-live/tmp/live-hourly.log; echo \"[\$(date '+%Y-%m-%d %H:%M:%S')] launchd trigger\" &gt;&gt; \$LOG_FILE; X_HEADLESS=1 /opt/homebrew/bin/node /Users/eric/.goldman-news-live/node_modules/tsx/dist/cli.mjs /Users/eric/.goldman-news-live/scripts/publish-live-feed.ts &gt;&gt; \$LOG_FILE 2&gt;&amp;1; STATUS=\$?; echo \"[\$(date '+%Y-%m-%d %H:%M:%S')] live:publish exit code: \$STATUS\" &gt;&gt; \$LOG_FILE; exit 0</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>HOME</key>
      <string>/Users/eric</string>
      <key>SHELL</key>
      <string>/bin/zsh</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>3600</integer>
  </dict>
</plist>
