# launchd/com.goldman.livefeed.hourly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.goldman.livefeed.hourly</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>'/Users/eric/Documents/Codex/Goldman news/11.1-----/scripts/run-live-publish.sh'</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/eric/Documents/Codex/Goldman news/11.1-----</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
      <key>HOME</key>
      <string>/Users/eric</string>
      <key>SHELL</key>
      <string>/bin/zsh</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>/Users/eric/Documents/Codex/Goldman news/11.1-----/tmp/live-launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eric/Documents/Codex/Goldman news/11.1-----/tmp/live-launchd.err.log</string>
  </dict>
</plist>
